# net-analyzer/ethtool-snmpd/files/ethtool-snmpd-initd-1
#!/sbin/openrc-run
# Copyright 2025 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

nameethtool-snmpd
pidfile=/run/ethtool-snmpd.pid
command=/usr/bin/ethtool-snmpd
command_args="${ETHTOOL_SNMPD_OPTS}"
