# net-analyzer/net-snmp/files/snmptrapd.init.2
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

SNMPTRAPD_PIDFILE="${SNMPTRAPD_PIDFILE:-/var/run/snmptrapd.pid}"

command="/usr/sbin/snmptrapd"
command_args="-p ${SNMPTRAPD_PIDFILE} ${SNMPTRAPD_FLAGS}"
pidfile="${SNMPTRAPD_PIDFILE}"

depend() {
	use logger
}
